# etc/cron.d/backup-dse
PATH=/usr/bin:/bin:/usr/sbin:/sbin
0 3 * * * ubuntu /usr/local/bin/backup-dse.sh